# install/packaging/armoreweb/debian/armoreweb.cron.d
#
# Regular cron jobs for the armoreweb package
#

20 *	* * *	root	 /usr/bin/armorewebmaint.sh
